# logs/danca2017_launch_20260515_182126/lure_rank0003_rank0005_missing_controls.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>APPLICATION_INSIGHTS_NO_STATSBEAT</key>
		<string>true</string>
		<key>BROWSER</key>
		<string>/Users/Xerk/.vscode-server/cli/servers/Stable-0958016b2af9f09bb4257e0df4a95e2f90590f9f/server/bin/helpers/browser.sh</string>
		<key>CODEX_CI</key>
		<string>1</string>
		<key>CODEX_INTERNAL_ORIGINATOR_OVERRIDE</key>
		<string>codex_vscode</string>
		<key>CODEX_SANDBOX_NETWORK_DISABLED</key>
		<string>1</string>
		<key>CODEX_THREAD_ID</key>
		<string>019e2dff-a7ab-7301-a98f-4b46018b43c6</string>
		<key>COLORTERM</key>
		<string></string>
		<key>DEBUG</key>
		<string>release</string>
		<key>ELECTRON_RUN_AS_NODE</key>
		<string>1</string>
		<key>GH_PAGER</key>
		<string>cat</string>
		<key>GIT_PAGER</key>
		<string>cat</string>
		<key>HIDDEN_ATTRACTORS_FRAC_ORDER</key>
		<string>0.9998</string>
		<key>HIDDEN_ATTRACTORS_MODEL</key>
		<string>piecewise</string>
		<key>HOME</key>
		<string>/Users/Xerk</string>
		<key>LANG</key>
		<string>C.UTF-8</string>
		<key>LC_ALL</key>
		<string>C.UTF-8</string>
		<key>LC_CTYPE</key>
		<string>C.UTF-8</string>
		<key>LOGNAME</key>
		<string>Xerk</string>
		<key>MPLCONFIGDIR</key>
		<string>/Users/Xerk/Desktop/Proyectos/Hidden Attractors Fractional Order/codigo_mac/.runtime_cache/matplotlib</string>
		<key>NO_COLOR</key>
		<string>1</string>
		<key>OLDPWD</key>
		<string>/Users/Xerk/Desktop/Proyectos/Hidden Attractors Fractional Order/codigo_mac</string>
		<key>OMP_NUM_THREADS</key>
		<string>1</string>
		<key>OMP_THREAD_LIMIT</key>
		<string>1</string>
		<key>PAGER</key>
		<string>cat</string>
		<key>PATH</key>
		<string>/Users/Xerk/.codex/tmp/arg0/codex-arg0cDzuwQ:/Users/Xerk/.vscode-server/cli/servers/Stable-0958016b2af9f09bb4257e0df4a95e2f90590f9f/server/bin/remote-cli:/Users/Xerk/Library/Python/3.10/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/usr/texbin:/opt/X11/bin:/Users/Xerk/.vscode-server/extensions/openai.chatgpt-26.513.21555-darwin-x64/bin/macos-x86_64</string>
		<key>PWD</key>
		<string>/Users/Xerk/Desktop/Proyectos/Hidden Attractors Fractional Order/codigo_mac</string>
		<key>PYTHONUNBUFFERED</key>
		<string>1</string>
		<key>RUST_LOG</key>
		<string>warn</string>
		<key>SHELL</key>
		<string>/bin/zsh</string>
		<key>SHLVL</key>
		<string>2</string>
		<key>SSH_CLIENT</key>
		<string>100.67.54.1 64249 22</string>
		<key>SSH_CONNECTION</key>
		<string>100.67.54.1 64249 100.76.77.5 22</string>
		<key>STM32CubeMX_PATH</key>
		<string>/Applications/STMicroelectronics/STM32CubeMX.app/Contents/Resources</string>
		<key>TERM</key>
		<string>dumb</string>
		<key>TMPDIR</key>
		<string>/var/folders/xf/6w3d0m1n5bx4_m1b7gfz9r3h0000gs/T/</string>
		<key>USER</key>
		<string>Xerk</string>
		<key>VSCODE_AGENT_FOLDER</key>
		<string>/Users/Xerk/.vscode-server</string>
		<key>VSCODE_CLI_REQUIRE_TOKEN</key>
		<string>5e802c65-7264-4656-ad87-70b1c3e8146d</string>
		<key>VSCODE_CWD</key>
		<string>/Users/Xerk</string>
		<key>VSCODE_ESM_ENTRYPOINT</key>
		<string>vs/workbench/api/node/extensionHostProcess</string>
		<key>VSCODE_HANDLES_SIGPIPE</key>
		<string>true</string>
		<key>VSCODE_HANDLES_UNCAUGHT_ERRORS</key>
		<string>true</string>
		<key>VSCODE_IPC_HOOK_CLI</key>
		<string>/var/folders/xf/6w3d0m1n5bx4_m1b7gfz9r3h0000gs/T/vscode-ipc-d446f922-8f5e-4152-8e97-d90a0d51f27e.sock</string>
		<key>VSCODE_NLS_CONFIG</key>
		<string>{"userLocale":"en","osLocale":"en","resolvedLanguage":"en","defaultMessagesFile":"/Users/Xerk/.vscode-server/cli/servers/Stable-0958016b2af9f09bb4257e0df4a95e2f90590f9f/server/out/nls.messages.json","locale":"en","availableLanguages":{}}</string>
		<key>VSCODE_RECONNECTION_GRACE_TIME</key>
		<string>10800000</string>
		<key>XDG_CACHE_HOME</key>
		<string>/Users/Xerk/Desktop/Proyectos/Hidden Attractors Fractional Order/codigo_mac/.runtime_cache/xdg_cache</string>
		<key>_</key>
		<string>/usr/local/bin/python3</string>
		<key>__CF_USER_TEXT_ENCODING</key>
		<string>0x1F9:0x0:0x0</string>
	</dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.hiddenattractors.danca2017.lure_missing.20260515_182126</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/python3</string>
		<string>/Users/Xerk/Desktop/Proyectos/Hidden Attractors Fractional Order/codigo_mac/lure_biased_multiparam_continuation.py</string>
		<string>--config</string>
		<string>/Users/Xerk/Desktop/Proyectos/Hidden Attractors Fractional Order/codigo_mac/logs/danca2017_launch_20260515_182126/lure_biased_multiparam_q09998_timestamped_config.json</string>
		<string>--post-continuation-only</string>
		<string>--survivor-id</string>
		<string>lure_biased_q_0p99980_rank_0003</string>
		<string>--survivor-id</string>
		<string>lure_biased_q_0p99980_rank_0005</string>
		<string>--execute-early-filter</string>
		<string>--execute-robustness</string>
		<string>--force</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/Xerk/Desktop/Proyectos/Hidden Attractors Fractional Order/codigo_mac/logs/danca2017_launch_20260515_182126/job_logs/lure_rank0003_rank0005_missing_controls.err</string>
	<key>StandardOutPath</key>
	<string>/Users/Xerk/Desktop/Proyectos/Hidden Attractors Fractional Order/codigo_mac/logs/danca2017_launch_20260515_182126/job_logs/lure_rank0003_rank0005_missing_controls.out</string>
	<key>WorkingDirectory</key>
	<string>/Users/Xerk/Desktop/Proyectos/Hidden Attractors Fractional Order/codigo_mac</string>
</dict>
</plist>
